# system_files/desktop/shared/usr/share/systemd/user/mute-audio-on-wake.service
[Unit]
Description=Mute audio after system wakes from suspend
After=suspend.target hibernate.target hybrid-sleep.target suspend-then-hibernate.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mute-audio-on-wake.sh

[Install]
WantedBy=suspend.target hibernate.target hybrid-sleep.target suspend-then-hibernate.target
